# supervisor/lara-app.conf
[program:lara-docker-starter]
command=php artisan queue:work redis --sleep=1 --tries=1
autostart=true
autorestart=true
user=root
numprocs=4
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/lara_docker_starter_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/lara_docker_starter_stdout.log
stdout_logfile_maxbytes=10MB
